# crontab
*/5 * * * * env DISPLAY=:0 /usr/bin/feh --randomize --bg-fill /home/mmillis/Pictures/bgs/*

